# .config/supervisor/local/supervisor-redis.conf
[program:redis]
command = /usr/local/bin/redis-server